# deploy/systemd/voiceasset-api.service
[Unit]
Description=VoiceAsset API
Requires=voiceasset-migrate.service
Wants=network-online.target
After=network-online.target voiceasset-migrate.service

[Service]
Type=simple
User=voiceasset
Group=voiceasset
WorkingDirectory=/data/apps/caddy/voice
EnvironmentFile=/data/apps/caddy/voice/config/server.env
ExecStart=/data/apps/caddy/voice/bin/voiceasset-api
Restart=on-failure
RestartSec=3s
TimeoutStopSec=20s
UMask=0027
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/data/apps/caddy/voice/objects
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
