# deploy/alertmanager/voiceasset-alert-receiver.service
[Unit]
Description=VoiceAsset local alert notification receiver
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=voiceasset
Group=voiceasset
WorkingDirectory=/data/apps/caddy/voice
ExecStart=/data/apps/caddy/voice/bin/voiceasset-alert-receiver --listen 127.0.0.1:19193 --output /data/apps/caddy/voice/alerts/notifications.jsonl
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/data/apps/caddy/voice/alerts
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
MemoryMax=64M
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
